# configurations/etc/init/mele-proxy.conf
# start mele-proxy
 
description     "MeleProxy"
 
start on virtual-filesystems 
stop on runlevel [06] 
 
expect fork
 
exec /usr/local/sbin/mele
