# configurations/etc/init/mele-proxy.conf
# start mele-proxy
 
description     "MeleProxy"
 
start on input-device-added ID_VENDOR_ID=1915
stop on runlevel [06] 
 
expect fork
 
exec /usr/local/sbin/mele
