# configurations/etc/init/mele-proxy.conf
# start mele-proxy
 
description     "MeleProxy"
 
start on desktop-session-start
stop on desktop-shutdown
 
expect fork
 
exec /usr/local/sbin/mele
